# compton/systemd/user/compton.service
[Unit]
Description=X compositor

[Service]
ExecStart=/usr/bin/compton

[Install]
WantedBy=multi-user.target
